# sonic/sonic.service
[Unit]
Description=Sonic Search Index
After=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/home/petar/.cargo/bin/sonic -c /etc/sonic.cfg
Restart=on-failure
# LimitNOFILE=4096

[Install]
WantedBy=multi-user.target